# src/deb/control/postinst
#!/bin/sh
sudo update-rc.d deltaservice defaults
